# frontend.ini
[program:content]
command=gunicorn editor_backend.wsgi:application --bind localhost:8080 --log-file /var/log/content_guni.log
directory=/var/frontend/editor_backend
autostart=true
stdout_logfile=/var/log/content.log
autorestart=true
startretries=3
