# frontend.ini
[program:content]
command=gunicorn editor_backend.wsgi:application --bind localhost:8080
directory=/var/content/editor_backend
autostart=true
autorestart=true
startretries=3
